# rootfs/etc/services.d/caddy/run
#!/usr/bin/execlineb -P
s6-setuidgid relicabackup
caddy -conf /etc/Caddyfile
