# rootfs/etc/services.d/caddy/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid daemon

caddy -conf /etc/Caddyfile
